# ota-agent/systemd/ota-agent.service
[Unit]
Description=OTA Agent
After=network.target

[Service]
ExecStart=/usr/bin/ota-agent
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
